# docker-elasticsearch/staging/etc/services.d/elasticsearch/finish
#!/usr/bin/with-contenv /bin/bash

source /opt/adp/bin/base-common.sh

$ADP_BIN/s6-terminate.sh
