# docker-logstash/staging/etc/services.d/logstash/run
#!/usr/bin/with-contenv /bin/bash

source /opt/adp/bin/base-common.sh

exec s6-setuidgid logstash $ADP_BIN/logstash-wrapper.sh
