# docker-logstash/staging/etc/services.d/gs-log-rest-server/run
#!/usr/bin/with-contenv /bin/bash

source /opt/adp/bin/base-common.sh

exec s6-setuidgid root $ADP_BIN/gs-log-rest-server.sh
